# flash-dumps/1552_build_2022_10_27_extracted/fs/usr/lib/opkg/info/busybox-init-base-files.list
/etc/openwrt_version
/etc/banner.failsafe
/var
/bin/setusbconfig
/etc/shells
/etc/device_info
/rdinit
/usr/bin/monitor-selfboot.sh
/usr/share/udhcpc/default.script
/etc/fstab
/etc/banner
/etc/init.d/rc.final
/etc/shadow
/etc/rc.common
/sbin/mount-sd.sh
/bin/adb_shell
/lib/functions/uci-defaults.sh
/etc/init.d/S41netparam
/lib/functions/network.sh
/etc/udev/rules.d/60-persistent-storage.rules
/etc/init.d/rc.preboot
/lib/functions/system.sh
/etc/mdev.conf
/pseudo_init
/etc/group
/etc/init.d/rc.modules
/lib/functions.sh
/lib/functions/leds.sh
/etc/init.d/rcK
/etc/init.d/rcS
/bin/adbd
/rom/note
/etc/inittab
/etc/udev/udev.conf
/etc/openwrt_release
/usr/bin/hotplug.sh
/etc/opkg/keys/af22f7a88858c8e9
/etc/profile
/lib/functions/service.sh
/etc/passwd
/etc/mtab
/lib/functions/preinit.sh
/etc/init.d/S10udev
