# flash-dumps/1552_build_2022_10_27_extracted/fs/usr/lib/opkg/info/swupdate.list
/usr/bin/progress
/etc/init.d/swupdate_autorun
/usr/bin/swupdate
/usr/bin/client
/sbin/swupdate_cmd.sh
