# flash-dumps/1552_build_2022_10_27_extracted/fs/etc/udev/rules.d/60-persistent-storage.rules
#/etc/udev/rules.d/60-persistent-storage.rules

# do not edit this file, it will be overwritten on update

# enable in-kernel media-presence polling
#ACTION=="add", SUBSYSTEM=="module", KERNEL=="block", ATTR{parameters/events_dfl_
#poll_msecs}=="0", ATTR{parameters/events_dfl_poll_msecs}="2000"
#ACTION=="add", ATTR{removable}=="1", ATTR{events_poll_msecs}=="-1", ATTR{events_
#poll_msecs}="2000"

SUBSYSTEM!="block", GOTO="persistent_storage_end"

ACTION=="add", KERNEL=="mmcblk[1-9]p[0-9]",OWNER="root",GROUP="root",MODE="0777"
,RUN+="/sbin/mount-sd.sh %k add"

ACTION=="remove", RUN+="/sbin/mount-sd.sh %k remove"

LABEL="persistent_storage_end"